# root/etc/services.d/tailscale-up/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

set -e

SOCKET="${TS_SOCKET:-/var/run/tailscale/tailscaled.sock}"

# Wait for tailscaled to be ready
echo "Waiting for tailscaled to be ready..."
while [[ ! -S "${SOCKET}" ]]; do
    sleep 1
done

# Build tailscale up command arguments
UP_ARGS=(
    "--socket=${SOCKET}"
    "up"
)

# Add authkey if provided
if [[ -n "${TS_AUTHKEY}" ]]; then
    UP_ARGS+=("--authkey=${TS_AUTHKEY}")
fi

# Add hostname if provided
if [[ -n "${TS_HOSTNAME}" ]]; then
    UP_ARGS+=("--hostname=${TS_HOSTNAME}")
fi

# Accept DNS configuration from tailnet
if [[ "${TS_ACCEPT_DNS}" == "true" ]]; then
    UP_ARGS+=("--accept-dns=true")
else
    UP_ARGS+=("--accept-dns=false")
fi

# Add any extra args for tailscale up
if [[ -n "${TS_EXTRA_ARGS}" ]]; then
    # shellcheck disable=SC2206
    UP_ARGS+=(${TS_EXTRA_ARGS})
fi

echo "Running tailscale up as user abc..."
s6-setuidgid abc tailscale "${UP_ARGS[@]}"

# Check current status
echo "Tailscale status:"
s6-setuidgid abc tailscale --socket="${SOCKET}" status

# Configure Taildrive shares if specified
# Format: "name1:path1,name2:path2"
if [[ -n "${TS_DRIVE_SHARES}" ]]; then
    echo "Configuring Taildrive shares..."
    IFS=',' read -ra SHARES <<< "${TS_DRIVE_SHARES}"
    for share in "${SHARES[@]}"; do
        name="${share%%:*}"
        path="${share#*:}"
        if [[ -n "${name}" && -n "${path}" ]]; then
            echo "  Adding share: ${name} -> ${path}"
            s6-setuidgid abc tailscale --socket="${SOCKET}" drive share "${name}" "${path}" || echo "  Warning: Failed to add share ${name} (is drive:share enabled in ACLs?)"
        fi
    done
    echo "Taildrive shares configured. List current shares:"
    s6-setuidgid abc tailscale --socket="${SOCKET}" drive list || true
fi

# This is a oneshot service - exit cleanly after bringing tailscale up
# s6 will not restart it unless we exit non-zero
echo "Tailscale is up. Service complete."
exec sleep infinity
