# root/etc/services.d/tailscaled/finish
#!/command/execlineb -S0
# shellcheck shell=bash

# Cleanup tailscaled on shutdown
foreground { tailscaled --cleanup }
foreground { redirfd -w 1 /run/s6-linux-init-container-results/exitcode echo 0 }
# /run/s6/basedir/bin/halt
